# docker/crontab
# Cron will run the artisan command every hour
0 * * * * root /usr/local/bin/php /app/artisan item:update
